# ops/relay-server/sharkbite-relay.service.example
[Unit]
Description=Sharkbite WebSocket relay
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/sharkbite/ops/relay-server
Environment=HOST=127.0.0.1
Environment=PORT=8082
ExecStart=/usr/bin/node node-server.mjs
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
